# client/sensos/stage-base/00-sensos/files/services/auto-hotspot.service
[Unit]
Description=WiFi Access Point and SSH Manager
After=network-online.target
Wants=network-online.target
Requires=NetworkManager.service

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'systemctl enable ssh && systemctl start ssh'
ExecStart=/bin/bash -c 'nmcli device wifi hotspot ssid sensosnet password sensossensos'
ExecStop=/bin/bash -c 'nmcli connection down Hotspot || true'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
