# client/sensos/stage-base/00-sensos/files/services/config-geekworm-eeprom.service
[Unit]
Description=Apply Geekworm UPS Configuration on First Boot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/config-geekworm-ups.sh
StandardError=journal
StandardOutput=journal
StandardInput=null

[Install]
WantedBy=multi-user.target
