# client/sensos/stage-base/00-sensos/files/services/ensure-sensos-dir.service
[Unit]
Description=Ensure /sensos directory structure exists with correct permissions
After=local-fs.target sensos-data.mount
RequiresMountsFor=/sensos/data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ensure-sensos-dir.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
