# client/sensos/stage-base/00-sensos/files/services/monitor-disk-space.service
[Unit]
Description=Disk Space Monitor and Service Controller

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitor-disk-space.sh

