# client/sensos/stage-base/00-sensos/files/services/send-status-update.service
[Unit]
Description=Send device status update (requires connectivity)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/requires-connectivity.sh /usr/local/bin/send-status-update.sh

