# client/sensos/stage-base/00-sensos/files/services/sensos-apply-bandwidth-policy.service
[Unit]
Description=Apply SensOS bandwidth policy (nftables lockdown)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/apply-bandwidth-policy.py
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
