# client/sensos/stage-base/00-sensos/files/services/sensos-apply-bandwidth-policy.path
[Unit]
Description=Watch /sensos/etc/network.conf and apply policy

[Path]
PathChanged=/sensos/etc/network.conf

[Install]
WantedBy=multi-user.target
