# client/sensos/stage-base/00-sensos/files/services/sensos-docker.service
[Unit]
Description=Start Docker Containers for Sound Analysis
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/environment

# Ensure we have a valid SENSOS_USER, then run the containers
ExecStart=/bin/bash -c '\
  if [ -n "$SENSOS_USER" ]; then \
    exec su - "$SENSOS_USER" -c "/usr/local/bin/start-sensos-containers.sh --no-detach"; \
  else \
    echo "SENSOS_USER is not set; skipping container startup." >&2; \
    exit 0; \
  fi'

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
